# crontab
SHELL=/bin/bash
30 0 * * * root /ocsp.sh >> /dev/null 2>&1
